# crontabs
* * * * * ~/Dropbox/config/bin/sticky-sync
